# files/default/oauth2_proxy.conf
description "oauth2 proxy"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
respawn
respawn limit 10 5
exec /usr/local/bin/oauth2_proxy --config /etc/oauth2_proxy.conf

